# 在线仓库/pegasus-master/vendor/topeet/file/topeet-start.service
#Copyright (c) 2025 Beijing Xunwei Electronics Co., Ltd
[Unit]
Description=TOPEET Hi3403V100 Start Ubuntu Script
 
[Service]
Type=oneshot
ExecStart=/etc/init.d/topeet-start.sh
RemainAfterExit=true
 
[Install]
WantedBy=sysinit.target
