# work/ss928_board_setup/systemd/ws73-bluetooth-module.service
[Unit]
Description=Load the SS928 WS73 Bluetooth controller module
After=local-fs.target
Before=bluetooth.service
ConditionPathExists=/opt/sample/ws73/ble_soc.ko

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ws73-bluetooth-module-start.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
