# work/linux_bmi270_backpack/bmi270-backpack.service
[Unit]
Description=BMI270 backpack posture monitor
After=bluetooth.service dbus.service
Wants=bluetooth.service dbus.service

[Service]
Type=simple
WorkingDirectory=/root/linux_bmi270_backpack
Environment=CONFIG=/root/linux_bmi270_backpack/config.ss928_ble.json
Environment=START_BLE_STACK=auto
ExecStart=/bin/sh -c 'exec /root/linux_bmi270_backpack/start_ss928_ble.sh >>/tmp/bmi270_ble.log 2>&1'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
